# src/main/config/server-test/root/.literacyapp/backup_prod/removeOldBackups.cron
# To add the cron job, run the following command: "crontab removeOldBackups.cron"

# Error recipient
MAILTO=info+crontab@elimu.ai

# minute hour day month weekday command
30 04 * * * /root/.literacyapp/backup_prod/removeOldBackups.sh

