# src/main/config/server-prod/root/.literacyapp/backup/mysqlBackup.cron
# To add the cron job, run the following command: "crontab mysqlBackup.cron"

# Error recipient
MAILTO=info+crontab@elimu.ai

# minute hour day month weekday command
00 04 * * * /root/.literacyapp/backup/mysqlBackup.sh
